# pi/scripts/99-hdmi-hotplug.rules
# HDMI hotplug udev rule - triggers resolution auto-detect when HDMI connected
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/local/bin/hdmi-hotplug.sh"
