# pi/mx5-display.service
[Unit]
Description=MX5 Telemetry Display
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
WorkingDirectory=/home/pi/MX5-Telemetry/pi/ui/src
# Kill any orphaned instances before starting (- prefix means ignore failures)
ExecStartPre=-/usr/bin/pkill -f 'python.*main.py'
# Wait for display/CAN to be ready
ExecStartPre=/bin/sleep 3
ExecStart=/usr/bin/python3 -u /home/pi/MX5-Telemetry/pi/ui/src/main.py --fullscreen
Restart=on-failure
RestartSec=10
# Ensure ALL child processes are killed on stop/restart
KillMode=control-group
# Timeout before SIGKILL
TimeoutStopSec=10

[Install]
WantedBy=graphical.target
